# rodiall-feed/rodial-ddns/files/etc/init.d/rodial-ddns
#!/bin/sh /etc/rc.common
START=99
STOP=10

start() {
    echo "rodial-ddns: one-shot update"
    /usr/bin/rodial-ddns || true
}
